# riemann_server/riemann.conf
#
# Upstart Ubuntu script for riemann
#

start on runlevel [2345]
stop on runlevel [06]

script
  exec start-stop-daemon --start --quiet --chuid root --exec '/home/pops/riemann-0.2.5/bin/riemann' '/home/pops/riemann-0.2.5/etc/riemann.config' 
end script

